# client/crontab
*/2 * * * * (cd /root/acri-olb/clinet; ruby ./olb-read.rb > ./olb-read.log)

