# client/crontab
*/2 * * * * ruby /root/acri-olb/client/olb-read.rb
