# client/crontab
*/2 * * * * (cd /root/acri-olb/client; ruby ./olb-read.rb > ./olb-read.log)

