# client/rc.local
#!/bin/sh

mkdir -p /tools
mount -t vboxsf tools -o ro,dmode=777 /tools

(cd /root/acri-olb/client; ruby olb-read.rb)
(cd /root/acri-olb/client; mv -f cleanfpga.log cleanfpga-shutdown.log; ./cleanfpga.sh)

