# server/crontab
*/30 * * * * /root/acri-olb/server/olb-monitor.rb
