# server/crontab
*/15 * * * * (cd /root/acri-olb/server; ./olb-monitor.rb > olb-monitor.log)
