# vm-host/crontab
00 */3 * * * /root/acri-olb/vm-host/restartvm.sh
