# vm-host/crontab
00 */6 * * * /root/acri-olb/vm-host/restartvm.sh
